# symfony/etc/services.d/nginx/finish
#!/bin/sh
# bring down container if nginx stops
s6-svscanctl -t /var/run/s6/services
